# datasource/ruuvitag-raspi/install/supervisor-proc.conf
[program:ruuvitag]
directory=/mnt/nfs/datasource/ruuvitag-raspi
command=npm start
autostart=true
autorestart=true
